# examples/systemd/project-backup.service
[Unit]
Description=Example Project encrypted backup
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/example-project-backup/backup.env
ExecStart=/opt/example-project/ops/backup/backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

